# debian/vote.service
[Unit]
Description=Vote Coloré - Real-time voting server
After=network.target
Documentation=https://github.com/robindubreuil/vote

[Service]
Type=simple
User=vote
Group=vote
EnvironmentFile=-/etc/default/vote
ExecStart=/usr/bin/vote-server
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
LockPersonality=true
ReadWritePaths=/var/lib/vote /run
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictRealtime=true
RestrictNamespaces=true
RestrictSUIDSGID=true
RemoveIPC=true
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
MemoryDenyWriteExecute=true

# Resource limits
LimitNOFILE=65536
TasksMax=512
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
