# debian/vote/usr/lib/systemd/system/vote.service
[Unit]
Description=Vote Coloré - Real-time voting server
After=network.target
Documentation=https://github.com/robin/vote

[Service]
Type=simple
User=vote
Group=vote
ExecStart=/usr/bin/vote-server
Restart=on-failure
RestartSec=5s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/vote /var/log/vote /run
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictRealtime=true
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Resource limits
LimitNOFILE=65536
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
